# deploy/curlycatclaw.service
[Unit]
Description=CurlyCatClaw personal AI assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=curlycatclaw
Group=curlycatclaw
ExecStart=/usr/local/bin/curlycatclaw --config /etc/curlycatclaw/config.toml
Restart=on-failure
RestartSec=5s

# Graceful shutdown: send SIGTERM, wait 35s (slightly more than the 30s internal timeout)
KillSignal=SIGTERM
TimeoutStopSec=35

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/curlycatclaw
# Uncomment if using file-based logging (logging.file in config.toml):
# ReadWritePaths=/var/log/curlycatclaw
ReadOnlyPaths=/etc/curlycatclaw
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
# Cannot set MemoryDenyWriteExecute=yes because wazero needs W^X for JIT.
# If wasm is disabled in config, you can set this to yes for tighter security.
MemoryDenyWriteExecute=no

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=curlycatclaw

[Install]
WantedBy=multi-user.target
